# crontab
* *   * * * root test -x /etc/squid3/update-squid-group-info.sh || /etc/squid3/update-squid-group-info.sh
