# crontab
* * * * * root /etc/squid3/update-squid-group-info.sh >> /var/log/squid3/update-squid-group-info.log 2>&1
